# config/supervisord.conf
[program:quant-autoresearch]
command=/Users/apple/quant-autoresearch/.venv/bin/python cli.py run --iterations 100 --safety high
directory=/Users/apple/quant-autoresearch
autostart=true
autorestart=true
startretries=3
stderr_logfile=/Users/apple/quant-autoresearch/run.err.log
stdout_logfile=/Users/apple/quant-autoresearch/run.out.log
user=apple
environment=GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",WANDB_API_KEY="%(ENV_WANDB_API_KEY)s"
